# workflows_dagster/systemd/dagster-webserver.service
[Unit]
Description=Dagster Webserver
Documentation=https://docs.dagster.io/
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=dagster
Group=dagster
WorkingDirectory=/opt/dagster-workflows/nebenkosten
Environment="DAGSTER_HOME=/opt/dagster-workflows/nebenkosten/workflows_dagster"
Environment="PYTHONPATH=/opt/dagster-workflows/nebenkosten:/opt/dagster-workflows/nebenkosten/Nebenkosten/src"
EnvironmentFile=/opt/dagster-workflows/nebenkosten/secrets/influxdb.env
EnvironmentFile=-/opt/dagster-workflows/nebenkosten/secrets/tibber.env
Environment="DAGSTER_POSTGRES_USER=dagster"
Environment="DAGSTER_POSTGRES_PASSWORD=dagster"
Environment="DAGSTER_POSTGRES_DB=dagster"
Environment="DAGSTER_POSTGRES_HOST=localhost"
Environment="DAGSTER_POSTGRES_PORT=5432"

# Use virtual environment Python
ExecStart=/opt/dagster-workflows/venv/bin/dagster-webserver -h 0.0.0.0 -p 3000

# Restart on failure
Restart=on-failure
RestartSec=10s

# Resource limits
MemoryMax=1G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dagster-webserver

[Install]
WantedBy=multi-user.target
